# deploy/Gentoo/funkwhale_server
#!/sbin/openrc-run

NAME=funkwhaleserver
PIDFILE=/var/run/$NAME.pid
USER=funkwhale
DAEMON_ARGS="-b 127.0.0.1 -p 5000 config.asgi:application --proxy-headers "
Daphne=/srv/funkwhale/virtualenv/bin/daphne
WORKDIR=/srv/funkwhale/api 

depend() {
        need net redis postgresql nginx funkwhale_beat funkwhale_worker
}

start() {
	ebegin "Starting Funkwhale Server"
        cd /srv/funkwhale/api
	set -a && source /srv/funkwhale/config/.env && set +a
        echo 'Starting Funkwhale Server'
	start-stop-daemon --start --user $USER --make-pidfile --pidfile $PIDFILE  -d $WORKDIR  --exec $Daphne -- $DAEMON_ARGS >> /var/log/funk/daphne.log 2>&1&
	echo 'Funkwhale Server started'
	echo
	eend $?
}

stop() {
	ebegin "Stopping Funkwhale"
	start-stop-daemon --stop --pidfile $PIDFILE
	eend $?
}
